# Suricata/crontabs
# 3am run
0 3 * * * /home/srogowski/rules/update_rules.sh